# openwrt/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/bin/sh /usr/bin/fix_sta_ap.sh > /dev/null &

exit 0